# luci-app-cowbping/root/etc/uci-defaults/luci-app-cowbping
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@cowbping[-1]
	add ucitrack cowbping
	set ucitrack.@cowbping[-1].init=cowbping
	commit ucitrack
EOF

touch /etc/config/cbp_cmd

rm -rf /tmp/luci-modulecache /tmp/luci-indexcache
exit 0
